# .config/systemd/user/rslsync.service
[Unit]
Description=Bittorent Sync per-user service
After=network.target

[Service]
Type=forking
ExecStart=/home/b268176/bin/rslsync
Restart=on-abort

[Install]
WantedBy=default.target
